# cookbooks/dmon/templates/default/dmon.conf.erb
description "dmon rest api service"
author      "Robert Cvitkovic"

start on filesystem or runlevel [2345]
stop on shutdown

setuid <%= @user %>
setgid <%= @group %>

exec <%= @install_dir %>/dmonEnv/bin/python \
  <%= @install_dir %>/src/start.py -p <%= @port %> > \
  <%= @install_dir %>/src/logs/dmon.log 2>&1
