# cookbooks/dmon/templates/default/dmon-es.conf.erb
description "ElasticSearch service"

start on (net-device-up
          and local-filesystems
          and runlevel [2345])
stop on runlevel [016]

respawn
respawn limit 10 30

limit nofile 32000 32000

setuid <%= @user %>
setgid <%= @group %>

env heap_size="<%= @heap_size %>"

script

  ES_HEAP_SIZE="$heap_size"

  export ES_HEAP_SIZE

  exec <%= @install_dir %>/bin/elasticsearch

end script
