# cookbooks/mongodb/templates/default/mongos.conf.erb
description "MongoDB shard router service"

limit nofile 20000 20000
kill timeout 60

start on started network-services
stop on runlevel [016]

setuid mongodb
setgid mongodb

exec /usr/bin/mongos --config /etc/mongod.conf
