# cookbooks/dmon_agent/templates/default/dmon_agent.conf.erb
description "dmon_agent rest api service"
author      "Robert Cvitkovic"

start on filesystem or runlevel [2345]
stop on shutdown

env DIR=<%= node['dmon_agent']['install_dir'] %>

setuid <%= node['dmon_agent']['user'] %>
setgid <%= node['dmon_agent']['group'] %>

script

    $DIR/dmonEnv/bin/python $DIR/dmon-agent.py > $DIR/log/dmon-agent.out 2>&1

end script
