# cookbooks/storm/templates/default/storm-daemon.conf.erb
description "Storm <%= @service %> upstart"

start on runlevel [2345]
stop on runlevel [016]

setuid <%= node['storm']['user'] %>
setgid <%= node['storm']['user'] %>

chdir <%= node['storm']['install_dir'] %>
exec bin/storm <%= @service %> 2>&1
respawn
