# cookbooks/dice-deployment-service/templates/default/celery.conf.erb
description "Celery worker"
author "tadeboro"

start on started uwsgi
stop on stopping uwsgi

respawn
respawn limit 20 5

setuid <%= @user %>
setgid <%= @user %>

exec <%= @venv %>/bin/celery worker \
   --app=dice_deploy \
   --queues=dice_deploy \
   --broker='amqp://guest@localhost//' \
   --loglevel=INFO \
   --concurrency=1 \
   --logfile=/var/log/celery/celery.log \
   --workdir=<%= @app_folder %> \
   --config=dice_deploy.settings \
   --without-gossip \
   --without-mingle \
   --events
