# cookbooks/dice-deployment-service/templates/default/uwsgi.conf.erb
description "uWSGI application server in Emperor mode"

start on runlevel [2345]
stop on runlevel [!2345]

setuid <%= @user %>
setgid www-data

exec <%= @app_venv %>/bin/uwsgi --logto /var/log/uwsgi/access.log \
  --emperor /etc/uwsgi/sites
