# cookbooks/cassandra/templates/default/cassandra.conf.erb
description "Cassandra upstart unit"

start on filesystem
stop on runlevel[016]

env CASSANDRA_HOME=<%= @install_dir %>
env CASSANDRA_CONF=/etc/cassandra

respawn
# if we have respawned once in the last 20 seconds then just stop.
respawn limit 2 20

setuid <%= @user %>
setgid <%= @group %>

console log

limit nofile 100000 100000
limit memlock unlimited unlimited

script
  <%= @install_dir %>/bin/cassandra -f > /dev/null 2>&1
end script
