# cookbooks/cassandra/templates/default/cassandra.conf.erb
description "Cassandra upstart unit"

start on filesystem
stop on runlevel[016]

env WAIT_FOR_START=20

env CASSANDRA_HOME=/usr/share/cassandra
env CASSANDRA_CONF=/etc/cassandra

respawn
# if we have respawned once in the last 20 seconds then just stop.
respawn limit 2 20

setuid <%= @user %>
setgid <%= @group %>

console log

chdir /

limit nofile 100000 100000
limit memlock unlimited unlimited

pre-start script

  [ -e /usr/share/cassandra/lib/apache-cassandra-<%= @version %>.jar ] || exit 0

end script

script
  /usr/share/cassandra/bin/cassandra -f > /dev/null 2>&1
end script

post-start script
  for tries in `seq $WAIT_FOR_START`; do
    sleep 1
    /usr/share/cassandra/bin/nodetool status && exit 0
  done
  exit 1
end script
