# cookbooks/cassandra/templates/default/cassandra.service.erb
[Unit]
Description=Cassandra
After=network.target

[Service]
User=<%= @user %>
Group=<%= @group %>
Environment="CASSANDRA_HOME=<%= @install_dir %>"
Environment="CASSANDRA_CONF=/etc/cassandra"
ExecStart=<%= @install_dir %>/bin/cassandra -f
StandardOutput=journal
StandardError=journal
LimitNOFILE=100000
LimitMEMLOCK=infinity
LimitNPROC=32768
LimitAS=infinity
Restart=always

[Install]
WantedBy=multi-user.target
